# k8s-installation/vagrant-ansible/ansible/roles/kubelet/templates/kubelet.service.j2
[Unit]
Description=Kubernetes Kubelet
Documentation=https://kubernetes.io/docs/
After=containerd.service
Requires=containerd.service

[Service]
ExecStart={{ kubelet_binary_dest }} \
    --node-name={{ node_name }} \
    --tls-cert-file={{ pki_dir }}/{{ kubelet_cert_file }} \
    --tls-private-key-file={{ pki_dir }}/{{ kubelet_key_file }} \
    {% for arg in kubelet_args %} {{ arg }} {% endfor %}
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target